# deploy/supervisord.railway.conf
[supervisord]
nodaemon=true
logfile=/dev/fd/1
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[program:worker]
directory=/srv/app
command=/bin/sh -lc 'exec celery -A app.workers:celery_app worker --loglevel=INFO --concurrency=${CELERY_CONCURRENCY:-2} --hostname=worker@%%h'
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:web]
directory=/srv/app
command=/bin/sh -lc 'exec uvicorn app.main:app --host 0.0.0.0 --port 8080'
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
